# systemd/waterflow.service
[Unit]
Description=Hydroponic Water Flow and Equipment Monitor
After=local-fs.target network-online.target systemd-udev-settle.service renogy.service
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__DIR__
EnvironmentFile=/etc/solar-hydroponic/credentials.env
ExecStart=/usr/bin/python3 __DIR__/waterflow_enhanced_failsafe.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
